# asset/overlay/etc/profile.d/qt-style.sh


#export QT_QPA_PLATFORMTHEME=kvantum
#export QT_PLATFORM_PLUGIN=kvantum
#export QT_STYLE_OVERRIDE=kvantum


export QT_QPA_PLATFORMTHEME=gtk2
export QT_PLATFORM_PLUGIN=gtk2
#export QT_STYLE_OVERRIDE=gtk2
